# ai_perception/face_recognition/edge_face_system/systemd/robot-face-attendance-logrotate.service
[Unit]
Description=Rotate Robot Face Attendance Logs

[Service]
Type=oneshot
User=gyul
Group=gyul
ExecStart=/home/gyul/robot_face_system/scripts/rotate_attendance_logs
