# ai_perception/face_recognition/edge_face_system/systemd/robot-face-attendance.service
[Unit]
Description=Robot Face Attendance Local Stack
Documentation=file:/home/gyul/robot_face_system/docs/current_operation_report.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=gyul
Group=gyul
WorkingDirectory=/home/gyul
Environment=CONFIG_PATH=/home/gyul/robot_face_system/config/edge_local.yaml
Environment=START_REALSENSE=true
Environment=LOG_DIR=/home/gyul/robot_face_system/logs
Environment=DATA_DIR=/home/gyul/robot_face_system/data
Environment=PYTHONUNBUFFERED=1
Environment=RCUTILS_LOGGING_MIN_SEVERITY=WARN
Environment=ROS_DOMAIN_ID=4
Environment=RMW_IMPLEMENTATION=rmw_fastrtps_cpp
Environment=ROS_AUTOMATIC_DISCOVERY_RANGE=SUBNET
Environment="ROS_STATIC_PEERS=192.168.0.48;"
ExecStart=/home/gyul/robot_face_system/scripts/run_attendance_local
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30
StandardOutput=append:/home/gyul/robot_face_system/logs/systemd_attendance.log
StandardError=append:/home/gyul/robot_face_system/logs/systemd_attendance.log

[Install]
WantedBy=multi-user.target
